# scripts/backup_mist.timer.example
[Unit]
Description=Daily Mist API backup

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true

[Install]
WantedBy=timers.target
